# heartbeat.service
[Unit]
Description=Daily monitoring heartbeat
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/heartbeat.sh