# mem-alert.timer
[Unit]
Description=Run mem-alert.sh every minute

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=10s

[Install]
WantedBy=timers.target