# mem-alert.service
[Unit]
Description=Memory pressure early-warning check
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mem-alert.sh