# tunnel-alert.timer
[Unit]
Description=Run tunnel-alert.sh every minute

[Timer]
OnBootSec=90
OnUnitActiveSec=60
AccuracySec=10s

[Install]
WantedBy=timers.target