# tunnel-alert.service
[Unit]
Description=Cloudflare Tunnel health watcher
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tunnel-alert.sh